# supervisor.conf
[inet_http_server]
port = 0.0.0.0:9001

[program:apache]
command = apachectl -DFOREGROUND

[program:feed_updater]
command = /srv/www/tt-rss/update_daemon2.php
startretries = 0
user = www-data
