# image/files/cron/crontab
@weekly root /usr/local/bin/update >> /var/log/tt-rss.log 2>&1
