# image/rootfs/etc/periodic/daily/update
#!/bin/sh
cd /srv/www/tt-rss
git pull
chown -R www-data:www-data .
echo $(date) TT-RSS was updated
