# image/rootfs/etc/services.d/feed-updater/run
#!/usr/bin/execlineb -P
s6-setuidgid www-data php8 /srv/www/tt-rss/update_daemon2.php
